# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:bot_ru]
command=python3 bot_ru.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/bot_ru.err.log
stdout_logfile=/var/log/bot_ru.out.log

[program:bot_hassp]
command=python3 bot_hassp.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/bot_hassp.err.log
stdout_logfile=/var/log/bot_hassp.out.log

[program:bot_smart]
command=python3 bot_smart.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/bot_smart.err.log
stdout_logfile=/var/log/bot_smart.out.log

[program:django]
command=gunicorn core.wsgi:application --bind 0.0.0.0:8000 --workers 2 --timeout 60
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/django.err.log
stdout_logfile=/var/log/django.out.log
startsecs=10
stopwaitsecs=10
environment=PYTHONUNBUFFERED="1"

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log
